# setup/pibox.service
[Unit]
Description=PiBox Autostart
Wants=network.target

[Service]
ExecStart=/bin/ bash autostart.sh > pibox.log
WorkingDirectory=/home/pi
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target